# hf-space/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:rag-api]
command=uvicorn src.api.main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app/rag
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    CHROMA_PERSIST_DIR="/app/rag/chroma_db",
    MLFLOW_TRACKING_URI="/app/rag/mlruns",
    LOG_FILE="/app/rag/logs/app.log",
    EMBEDDING_DEVICE="cpu",
    API_RELOAD="false"

[program:cv-api]
command=uvicorn src.api.main:app --host 0.0.0.0 --port 8001 --workers 1
directory=/app/cv
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    CV_DEVICE="cpu",
    CV_MODELS_CACHE="/app/cv/model_cache",
    CV_UPLOAD_DIR="/app/cv/uploads",
    MLFLOW_TRACKING_URI="/app/cv/mlruns",
    CV_API_PORT="8001"
